# supervisord.ini
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=root:root

[supervisord]
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:django]
directory = /app/backend
command = python manage.py runserver 0.0.0.0:8000
autostart = false
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:react]
directory=/app/frontend
command=npm start
autostart=false
autorestart=true
stdout_logfile=supervisord.log
stderr_logfile=supervisord.err